# container/supervisor/conf.d/adminrouter.conf
[program:adminrouter]
command=nginx -p /etc/nginx -c /etc/nginx/conf/nginx.master.conf -g 'daemon off;'
startretries=99
