# container/supervisor/conf.d/mesos-dns.conf
[program:mesos-dns]
command=mesos-dns -config=/etc/mesos-dns-config.json
